# archiso/airootfs/etc/profile.d/endeavour-env.sh
#!/bin/bash
# /etc/profile.d/endeavour-env.sh
# Environment variables for Endeavour OS Privacy Edition

export ENDEAVOUR_VERSION="2024.1"
export ENDEAVOUR_NAME="Endeavour OS Privacy Edition"

echo "Welcome to $ENDEAVOUR_NAME"
